# debian/cron.daily
#!/bin/sh

. /etc/default/mythtv-status

if [ "x$EMAIL" = "x" -o "x$EMAIL" = "xnone" ]
then
  exit
fi

/usr/bin/mythtv-status --host $HOST --email $EMAIL $EMAIL_ARGS

# Keep run-parts happy by always returning with a return code of 0.
exit 0
