# debian/cron.d
#
# Regular cron jobs for the mythtv-status package
#
*/10 *	* * *	root	/etc/init.d/mythtv-status reload > /dev/null
